# templates/docker-extra.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Docker extra: {{ item.key }}
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true

Environment="TERM=dumb"

WorkingDirectory=/opt/docker-extra_{{ item.key }}

ExecStart=/usr/bin/sh -c 'docker ps --format "{% raw %}{{.Names}}{% endraw %}" -a | grep -q -x {{ item.key }} && docker start {{ item.key }} || docker run -d --name {{ item.key }} --env-file .env {{ item.value.parameters }}'
ExecStop=/usr/bin/sh -c 'docker stop {{ item.key }}; docker rm -f {{ item.key }}'

[Install]
WantedBy=multi-user.target
